# systemd/flow2api-host-agent-daily-restart.service
[Unit]
Description=Flow2API Host Agent - daily maintenance restart

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart flow2api-host-agent-browser.service flow2api-host-agent.service flow2api-host-agent-ui.service
